# buhayra/jobs/crontab

SHELL=/bin/bash
MAILTO="jmigueldelgado@gmail.com"


### acquire scenes from ESA

0 20 * * 1 /users/stud09/martinsd/local/miniconda2/envs/snappy/bin/python2 /users/stud09/martinsd/proj/buhayra "getscenes"
0 20 * * 5 /users/stud09/martinsd/local/miniconda2/envs/snappy/bin/python2 /users/stud09/martinsd/proj/buhayra "getscenes"

### process acquired scenes
0 2 * * * qsub /users/stud09/martinsd/proj/buhayra/buhayra/jobs/2_sar2watermask.job
#0 2 * * * qsub /users/stud09/martinsd/proj/buhayra/buhayra/jobs/2.1_doCloudMask.job
0 7 * * * qsub /users/stud09/martinsd/proj/buhayra/buhayra/jobs/2.2_n2w.job

### polygonize scenes from sentinel 1 and 2
0 14 * * * qsub /users/stud09/martinsd/proj/buhayra/buhayra/jobs/3_polygonize.job

### simplify polygons and save as geojson
0 13 * * * qsub /users/stud09/martinsd/proj/buhayra/buhayra/jobs/4_st_simplify.job

### insert polygons into mongodb database in webserver
0 23 * * 3 /users/stud09/martinsd/proj/buhayra/buhayra/jobs/insert_polys_mongo.sh
0 23 * * 5 /users/stud09/martinsd/proj/buhayra/buhayra/jobs/insert_polys_mongo.sh
0 23 * * 7 /users/stud09/martinsd/proj/buhayra/buhayra/jobs/insert_polys_mongo.sh

### clean up
#52 15 * * 5 /users/stud09/martinsd/local/miniconda2/envs/mongodb/bin/python2 /users/stud09/martinsd/proj/buhayra/buhayra/cleanup.py


# download scenes mondays and wednesdays at 8 pm
# process data every day at 2 am to retry in case wall time is violated
# polygonize every day at 2 pm
# simplify in R and save as geojson every day at 1 pm
# insert json documents into mongodb in webserver wednesdays, fridays and sundays at 11 pm
